# dashboard-mvp/infra/systemd/backup@prune.service
[Unit]
Description=ClickHouse Backup Prune
Documentation=man:clickhouse-backup(1)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'cd /opt/zakaz_dashboard && ./ops/backup_prune.sh'
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600
TimeoutStopSec=300

# Переменные окружения
EnvironmentFile=-/opt/zakaz_dashboard/.env

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/clickhouse/backups /tmp

# Перезапуск при сбое
Restart=on-failure
RestartSec=300
StartLimitBurst=2
StartLimitInterval=86400

[Install]
WantedBy=multi-user.target