# dashboard-mvp/infra/systemd/etl@.service
[Unit]
Description=ETL Job %i
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/dashboard-mvp
EnvironmentFile=/opt/dashboard-mvp/.env
ExecStart=/usr/bin/bash -lc "ops/run_job.sh %i"
# логи в journald
TimeoutStartSec=900