# dashboard-mvp/ops/systemd/qtickets_sheets.service
[Unit]
Description=QTickets Google Sheets Data Loader
After=network.target clickhouse.service
Wants=clickhouse.service

[Service]
Type=oneshot
User=etl
Group=etl
WorkingDirectory=/opt/zakaz_dashboard
Environment=PYTHONPATH=/opt/zakaz_dashboard
ExecStart=/usr/bin/python3 -m integrations.qtickets_sheets.loader --envfile secrets/.env.qtickets_sheets --ch-env secrets/.env.ch --verbose
EnvironmentFile=secrets/.env.qtickets_sheets
EnvironmentFile=secrets/.env.ch
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qtickets-sheets-loader
Restart=on-failure
RestartSec=30
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target